# src/supervisord.conf
[supervisord]
nodaemon=true

[program:sandstorm]
command=/opt/sandstorm/sandstorm/sandstorm start

